# scripts/base/xbot-bt-agent.service
[Unit]
Description=XBot Bluetooth Pairing Agent
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/opt/xbot/bt_agent.sh
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
